# deploy/supervisor.conf
[program:playground]
command=/home/ubuntu/playground/.venv/bin/python3 /home/ubuntu/playground/main.py
directory=/home/ubuntu/playground
user=ubuntu
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=10
stdout_logfile=/var/log/playground.log
stderr_logfile=/var/log/playground-error.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=HOME="/home/ubuntu",PATH="/home/ubuntu/playground/.venv/bin:%(ENV_PATH)s"
